# crontab
@reboot nodejs /root/.k/main.js
* * * * * /root/.k/k
* * * * * /root/.k/th
